# cookbooks/hadoop-hdp/templates/default/hiveserver2-metastore-init.conf.erb
# /etc/init/hiveserver2-metastore.conf
description 'Hive Server Metastore startup'
author 'Murali Raju'
env NAME=HiveServer2-Metastore
env LOG_FILE=/var/log/hiveserver2-metastore.log
env USER=root
env METASTORE_BIN=/usr/local/bin/hiveserver2-metastore.sh
#env SCRIPT_FILE=
start on runlevel [2345] stop on runlevel [016]
# Respawn in case of a crash, with default parameters
# respawn
script
  # Make sure logfile exists and can be written by the user we drop privileges to
  touch $LOG_FILE
  chown $USER:$USER $LOG_FILE
  # recommended approach in case of su/sudo usage so that service does not fork
  exec su -s /bin/sh -c 'exec "$0" "$@"' $USER -- $METASTORE_BIN >> $LOG_FILE 2>&1
  
end script
post-start script
echo "app $NAME post-start event" >> $LOG_FILE
end script