# cookbooks/hadoop-hdp/templates/default/oozie-init.conf.erb
# /etc/init/oozie-server.conf
description 'Oozie Server startup'
author 'Murali Raju'
env NAME=OozieServer
env LOG_FILE=/var/log/oozie-server.log
env USER=oozie
env SCRIPT_FILE=/usr/local/bin/start-oozie.sh
start on runlevel [2345] stop on runlevel [016]
# Respawn in case of a crash, with default parameters
# respawn
script
  # Make sure logfile exists and can be written by the user we drop privileges to
  touch $LOG_FILE
  chown $USER:$USER $LOG_FILE
  # recommended approach in case of su/sudo usage so that service does not fork
  exec su -s /bin/sh -c 'exec "$0" "$@"' $USER -- $SCRIPT_FILE >> $LOG_FILE 2>&1
end script
post-start script
echo "app $NAME post-start event" >> $LOG_FILE
end script